# configs/supervisor.d/duplicati.ini
[program:duplicati]
command=/usr/bin/mono /opt/duplicati/Duplicati.Server.exe --webservice-port=8200 --log-file=/var/log/duplicati.log --log-level=Information --log-retention=1m --webservice-password=__PASSWORD__ --server-datafolder=/data/ --webservice-interface=any
redirect_stderr=true
stdout_logfile=/var/log/supervisor/duplicati.log
